# docker/production/supervisord-meteor.conf
[program:meteor]
directory = /code/bundle/
command = node main.js
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /log/supervisor-meteor.log
