# docker/production/supervisord-thrift.conf
[program:thrift]
directory = /thrift
command = /thrift/LinkedHashtags_thrift_RPC.out
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /log/supervisor-thrift.log
