# docker/production/supervisord-mongodb.conf
[program:mongodb]
command = mongod -dbpath=/data
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /log/supervisor-mongodb.log
